# launchd/com.dansbytracker.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.dansbytracker.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd /Users/ai_lab/dansbytracker &amp;&amp; source venv/bin/activate &amp;&amp; SENDER_EMAIL=labairj@gmail.com SAVANT_OAA_CSV_URL=/Users/ai_lab/dansbytracker/outs_above_average.csv SAVANT_FRV_CSV_URL=/Users/ai_lab/dansbytracker/fielding_run_value.csv python3 dansbytracker.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/ai_lab/dansbytracker</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/ai_lab/dansbytracker/launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ai_lab/dansbytracker/launchd.log</string>

</dict>
</plist>
